# packaging/systemd/meshcore-packet-capture.service
[Unit]
Description=MeshCore Packet Capture
After=time-sync.target network-online.target bluetooth.target
Wants=time-sync.target network-online.target

[Service]
Type=exec
User=meshcore-capture
Group=meshcore-capture
WorkingDirectory=/opt/meshcore-packet-capture
ExecStart=/opt/meshcore-packet-capture/venv/bin/python3 -m meshcore_packet_capture
# Best-effort BLE teardown after the process exits (no-op for serial/TCP;
# the script itself gates on connection_type="ble" + a configured ble_address).
ExecStopPost=/opt/meshcore-packet-capture/ble-disconnect.sh
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/meshcore-packet-capture /var/lib/meshcore-packet-capture

[Install]
WantedBy=multi-user.target
